# deploy/fab/supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".

[unix_http_server]
file=/home/ubuntu/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/home/ubuntu/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=debug                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/ubuntu/supervisor.sock ; use a unix:// URL  for a unix socket

; Regardig stopasgroup=true - see:
;   * https://gist.github.com/badsyntax/0e937159f26ac82f41fa
;   * https://coderwall.com/p/4tcw7w/setting-supervisor-to-really-stop-django-runserver

[program:flask]
command="taxbrain-flask-worker"
user=ubuntu
environment=PATH="/home/ubuntu/deploy/:/home/ubuntu/miniconda2/envs/aei_dropq/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games"
directory=/home/ubuntu/deploy
stopasgroup=true


[program:celery]
command="taxbrain-celery-worker"
user=ubuntu
environment=PATH="/home/ubuntu/deploy/:/home/ubuntu/miniconda2/envs/aei_dropq/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games"
directory=/home/ubuntu/deploy
stopasgroup=true


