# distributed/api/supervisord_local.conf
[supervisorctl]
username=cloud-user

[supervisord]
stdout_logfile=/dev/fd/1;
stdout_logfile_maxbytes=0;
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[program:flask]
command=python flask_server.py
priority=1
killasgroup=true

[program:celery]
command=celery -A celery_tasks worker --concurrency=1 -P eventlet -l info
priority=1
killasgroup=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
